# scripts/plists/com.trade.update-all.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>Label</key>
	<string>com.trade.update-all</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/linhuichen/code/trade/scripts/update_all.sh</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/linhuichen/code/trade/data/logs/update_all_launchd.err</string>
	<key>StandardOutPath</key>
	<string>/Users/linhuichen/code/trade/data/logs/update_all_launchd.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>17</integer>
		<key>Minute</key>
		<integer>50</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/linhuichen/code/trade</string>
	<!-- 超时 7200s：全量采集最重（多源 + 回填 + 推送），给足上限防僵死。 -->
	<key>ExitTimeOut</key>
	<integer>7200</integer>
</dict>
</plist>
