# scripts/plists/com.trade.rzhb-backfill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trade.rzhb-backfill</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/linhuichen/code/trade/scripts/rzhb_backfill.sh</string>
    </array>
    <!-- 每日 23:00 跑：两融余额（a_fund_margin，沪市融资余额）源 stock_margin_sse
         盘后发布当日数据，实测 22:10 时 latest 仍=前一日（源发布晚于22:10）。
         23:00 当晚单采（比22:10晚50分钟更可能采到当日），采到则当日上线；
         配合凌晨 backfill 缺口检测兜底（commit 61567f0）治本——源晚出也能补。
         独立于 20:00 backfill 的 SERIES_FUNCS 兜底；采集几秒 + compute + 持 deploy 锁推送。 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/linhuichen/code/trade</string>
    <key>StandardOutPath</key>
    <string>/Users/linhuichen/code/trade/data/logs/rzhb_backfill_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/linhuichen/code/trade/data/logs/rzhb_backfill_launchd.err</string>
    <key>RunAtLoad</key>
    <false/>
    <!-- 超时 3600s：deploy.lock 可能排队等晚跑任务，给足等待。 -->
    <key>ExitTimeOut</key>
    <integer>3600</integer>
</dict>
</plist>
