# scripts/plists/com.trade.futures-backfill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trade.futures-backfill</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/linhuichen/code/trade/scripts/futures_backfill.sh</string>
    </array>
    <!-- 每日 20:05(主槽) + 21:00(兜底二次槽) 跑：CFFEX 股指期货前20名会员持仓次日约 20:00 后发布当日数据，
         17:50 update_all 采不到，20:05 单采当日根治 46h 滞后。
         采集几秒 + 持 deploy 锁推送；20:00 backfill / 20:07 etf 密集时段 deploy.lock 排队。
         21:00 二次槽：应对 20:05 网络抖动/服务器异常(CFFEX~20:00出，21:00大概率采到)，幂等覆盖重采。 -->
    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>5</integer>
      </dict>
      <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/linhuichen/code/trade</string>
    <key>StandardOutPath</key>
    <string>/Users/linhuichen/code/trade/data/logs/futures_backfill_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/linhuichen/code/trade/data/logs/futures_backfill_launchd.err</string>
    <key>RunAtLoad</key>
    <false/>
    <!-- 超时 3600s：deploy.lock 可能排队等 20:00 backfill，给足等待。 -->
    <key>ExitTimeOut</key>
    <integer>3600</integer>
</dict>
</plist>
