# launchd/com.trade.sentiment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trade.sentiment</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/linhuichen/code/trade/.venv/bin/python</string>
        <string>-m</string>
        <string>app.scheduler</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>15</integer>
        <key>Minute</key>
        <integer>33</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/linhuichen/code/trade</string>
    <key>StandardOutPath</key>
    <string>/Users/linhuichen/code/trade/data/logs/scheduler.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/linhuichen/code/trade/data/logs/scheduler.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
